# crontab
*/15 * * * * echo 'EXPIRE START' && sillypaste expire && echo 'EXPIRE END'
